# launchd/com.mtplx.context-sync.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mtplx.context-sync</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HERMES_INFRA_DIR</key>
        <string>/path/to/hermes-infra</string>
    </dict>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/hermes-infra/target/release/hermes-infra</string>
        <string>mtplx-context-sync</string>
    </array>
    <key>StartInterval</key>
    <integer>2</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/path/to/hermes-infra</string>
    <key>StandardOutPath</key>
    <string>/path/to/hermes-infra/logs/context-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/hermes-infra/logs/context-sync.log</string>
</dict>
</plist>
