# launchd/com.hermes.tei.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.tei</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HERMES_INFRA_DIR</key>
        <string>/path/to/hermes-infra</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/hermes-infra/target/release/hermes-infra</string>
        <string>tei</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>ExitTimeOut</key>
    <integer>15</integer>

    <!-- Replace $HOME with your actual home directory path -->
    <key>StandardOutPath</key>
    <string>/path/to/hermes-infra/logs/tei.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/hermes-infra/logs/tei.log</string>

    <!-- TEI needs ~45s to warm up Metal on Apple Silicon before serving requests.
         Both Hindsight and the code-indexer depend on this endpoint. -->
    <key>LaunchOnlyOnce</key>
    <false/>
</dict>
</plist>
