# launchd/com.hermes.hindsight.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HERMES_INFRA_DIR</key>
        <string>/path/to/hermes-infra</string>
        <key>HINDSIGHT_API_EMBEDDINGS_DIM</key>
        <string>1024</string>
        <key>HINDSIGHT_API_EMBEDDINGS_OPENAI_API_KEY</key>
        <string>tei</string>
        <key>HINDSIGHT_API_EMBEDDINGS_OPENAI_BASE_URL</key>
        <string>http://localhost:6999/v1</string>
        <key>HINDSIGHT_API_EMBEDDINGS_OPENAI_MODEL</key>
        <string>Qwen/Qwen3-Embedding-0.6B</string>
        <key>HINDSIGHT_API_EMBEDDINGS_PROVIDER</key>
        <string>openai</string>
        <key>HINDSIGHT_API_RERANKER_PROVIDER</key>
        <string>rrf</string>
        <key>HINDSIGHT_PROFILE</key>
        <string>hermes</string>
    </dict>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>Label</key>
    <string>com.hermes.hindsight</string>
    <key>LaunchOnlyOnce</key>
    <false/>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/hermes-infra/target/release/hermes-infra</string>
        <string>hindsight</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/path/to/hermes-infra/logs/hindsight-api.log</string>
    <key>StandardOutPath</key>
    <string>/path/to/hermes-infra/logs/hindsight-api.log</string>
    <key>WorkingDirectory</key>
    <string>/path/to/.hermes/hindsight</string>
</dict>
</plist>
